# surface/systemd/infoscreen-live-data.timer
[Unit]
Description=Run InfoScreen live data fetch every 5 minutes

[Timer]
OnBootSec=30
OnUnitActiveSec=5min
Unit=infoscreen-live-data.service

[Install]
WantedBy=timers.target
