# deploy/systemd/user/infoscreen-event-stream.service
[Unit]
Description=Fetch RSS event stream for InfoScreen

[Service]
Type=oneshot
ExecStart=/usr/bin/python3 %h/infoscreen/fetch_event_stream.py
