# deploy/systemd/user/infoscreen-http.service
[Unit]
Description=InfoScreen HTTP server

[Service]
WorkingDirectory=%h/infoscreen/surface
ExecStart=/usr/bin/python3 %h/infoscreen/surface/serve_infoscreen.py
Restart=always
RestartSec=2
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
